# modules/xvfb/templates/xvfb.service.erb
[Unit]
Description=X Virtual Frame Buffer Service
After=network.target

[Service]
ExecStart=/usr/bin/Xvfb <%= @display %> -screen 0 <%= @resolution %>

[Install]
WantedBy=multi-user.target
